# sr7400.conf
# /etc/init/sr7400.conf
# upstart script. Runs the app as a daemon, including manaing PID files etc
# usage:
#   sudo start|stop|status sr7400

# Output to /var/log/upstart/sr7400.log
console log

description "Upstart script to run an app as a service or daemon"
author "Andrew Cuddon"
version "1.0"

env VM=/usr/local/bin/node
env APP_DIR=/home/ubuntu/avc3/sr7400
env SCRIPT_FILE="ps3monitor.js"
# Start as root and drop privileges to user:
env USER="ubuntu"
env GROUP="ubuntu"

# Start service on startup, stop on shutdown
start on (filesystem and net-device-up IFACE=eth0)
stop on runlevel [!2345]

# Automatically restart process if crashed
respawn

pre-start script
    DATETIME=`/bin/date +%D-%I:%M:%S`
    echo "===== Upstart: Preparing script: $UPSTART_JOB: $DATETIME ====="
    # Remove an existing log file
    if [ -f /var/log/upstart/$UPSTART_JOB.log ]; then
        rm /var/log/upstart/$UPSTART_JOB.log
        echo "===== Upstart: Preparing script: $UPSTART_JOB: $DATETIME ====="
        echo "Existing logile deleted"
    fi
end script

# Main process to daemonise
exec start-stop-daemon --start --make-pidfile --pidfile /var/run/$UPSTART_JOB.pid --chuid $USER:$GROUP --chdir $APP_DIR --exec $VM $SCRIPT_FILE


post-start script
    # Script or process to run after the main process has been spawned, but before the started(7) event has been emitted
    DATETIME=`/bin/date +%D-%I:%M:%S`
    echo "===== Upstart: $UPSTART_JOB started: $DATETIME ====="
end script

post-stop script
    DATETIME=`/bin/date +%D-%I:%M:%S`
    echo "===== Upstart: $UPSTART_JOB stopped: $DATETIME ====="
end script


